# supervisor.conf
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:app]
command=serverless offline start

[program:test]
command=./node_modules/.bin/jest